# web/nginx + guincorn/6/2/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:gunicorn]
command=gunicorn main:app -c gunicorn.py
directory=/deploy