# web/nginx + guincorn/7/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:gunicorn]
command=gunicorn -c /app/gunicorn.py HelloWorld.wsgi:application