# web/nginx + guincorn/2/gunicorn.conf
[program:gunicorn]
command=gunicorn app:app -b localhost:5000
directory=/deploy