# web/nginx/1/tomcats/supervisord.conf
[supervisord]
nodaemon=true

[program:main1]
command=python /src/main1.py

[program:main2]
command=python /src/main2.py

[program:nginx]
command=nginx