# web/nginx + guincorn/6/1/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

[program:gunicorn]
command=gunicorn app:app -c gunicorn.py
directory=/deploy